# etc/systemd/system/monerod-prune.service
[Unit]
Description=Monero Updater
After=network.target

[Service]
ExecStart=/bin/bash monerod-prune.sh
WorkingDirectory=/home/nanode/
StandardOutput=inherit
StandardError=inherit
Type=oneshot
RemainAfterExit=true
Restart=no
User=nanode

[Install]
WantedBy=multi-user.target
