# etc/systemd/system/monerod-start-free.service
[Unit]
Description=Monero Free Public Node
After=network.target

[Service]
ExecStart=/bin/bash monerod-start-public-free.sh
WorkingDirectory=/home/pinodexmr/
Type=forking
Restart=always
RestartSec=30
User=pinodexmr
StandardOutput=inherit
StandardError=inherit
PIDFile=/home/pinodexmr/monero/build/release/bin/monerod.pid

[Install]
WantedBy=multi-user.target

